# scripts/installation/rc.local
...
rfcomm bind 0 [MAC ADDRES]
systemctl start smarthome
...
